# deploy/supervisor.conf
; ============================================
; Supervisor config for Article Factory
; ============================================
; Install: sudo apt install supervisor
; Copy to: /etc/supervisor/conf.d/article-factory.conf
; Apply: sudo supervisorctl reread && sudo supervisorctl update

[program:article-worker]
command=/opt/article-factory/venv/bin/python /opt/article-factory/worker.py
directory=/opt/article-factory
user=www-data
process_name=%(program_name)s_%(process_num)02d
numprocs=2
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
stdout_logfile=/opt/article-factory/logs/worker-%(process_num)02d.log
stderr_logfile=/opt/article-factory/logs/worker-%(process_num)02d-error.log
environment=PATH="/opt/article-factory/venv/bin:%(ENV_PATH)s"

[program:rq-dashboard]
command=/opt/article-factory/venv/bin/rq-dashboard --port 9181 --redis-url redis://localhost:6379/0
directory=/opt/article-factory
user=www-data
autostart=true
autorestart=true
stdout_logfile=/opt/article-factory/logs/dashboard.log
stderr_logfile=/opt/article-factory/logs/dashboard-error.log
